# backup.timer
[Unit]
Description=melee.cool backup

[Timer]
OnCalendar=weekly
AccuracySec=12h
Persistent=true

