# services/certbot.timer
[Unit]
Description=Certbot Renewal Timer

[Timer]
OnCalendar=monthly
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
